# supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001

[supervisord]
logfile=/var/log/ingest-supervisord.log
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
serverurl=http://127.0.0.1:9001

[program:mongo]
command=mongod --config mongod.conf --rest
process_name=%(program_name)s 
stdout_logfile=/var/log/ftd-mongod.log
stderr_logfile=/var/log/ftd-mongod.err.log
priority=1

[program:rabbitmq]
command=/usr/local/sbin/rabbitmq-server
process_name=%(program_name)s 
stdout_logfile=/var/log/tremble-%(program_name)s.log
stderr_logfile=/var/log/tremble-%(program_name)s.err.log
priority=1                    
